# node/crontab
0 * * * * rsync -av /psv_results/ /mnt/nas-backup/ethoscope_results  --exclude=lost+found --exclude="*.db-journal" --exclude="*.zip"
